# crontab
17 4 * * * /dir-to-b2/backup >> /var/log/cron.log 2>&1
